# rootfs/etc/systemd/system/first-boot.service
[Unit]
Description=DCS first-boot enrollment
After=network-online.target
Wants=network-online.target
ConditionPathExists=!/var/lib/dcs/enrolled

[Service]
Type=oneshot
ExecStart=/usr/local/sbin/dcs-enroll
StandardOutput=journal
StandardError=journal
Restart=on-failure
RestartSec=30s

[Install]
WantedBy=multi-user.target
