# rootfs/etc/systemd/system/first-boot.service
[Unit]
Description=Detel first-boot enrollment
After=network-online.target
Wants=network-online.target
ConditionPathExists=!/var/lib/detel/enrolled

[Service]
Type=oneshot
ExecStart=/usr/local/sbin/detel-enroll
StandardOutput=journal
StandardError=journal
Restart=on-failure
RestartSec=30s

[Install]
WantedBy=multi-user.target
