# rootfs/etc/systemd/system/subterra-heartbeat.service
[Unit]
Description=Subterra tunnel heartbeat
After=wg-quick@wg0.service
Wants=wg-quick@wg0.service

[Service]
Type=oneshot
ExecStart=/usr/local/sbin/subterra-heartbeat
